# packages/cli/scripts/service/com.chainlesschain.bridge-mtc.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd LaunchAgent for `cc crosschain mtc-serve` on macOS.

  Install:
    cp com.chainlesschain.bridge-mtc.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.chainlesschain.bridge-mtc.plist

  Stop:
    launchctl unload ~/Library/LaunchAgents/com.chainlesschain.bridge-mtc.plist

  Logs land in ~/Library/Logs/cc-bridge-mtc/{out,err}.log
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.chainlesschain.bridge-mtc</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/cc</string>
        <string>crosschain</string>
        <string>mtc-serve</string>
        <string>--interval</string>
        <string>60</string>
        <string>--json</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>5</integer>

    <key>StandardOutPath</key>
    <string>/Users/YOUR_USER/Library/Logs/cc-bridge-mtc/out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USER/Library/Logs/cc-bridge-mtc/err.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/YOUR_USER</string>
</dict>
</plist>
